# .docker/fpm/var/spool/cron/crontabs/www-data
*/10 * * * * php /var/www/bin/satis build /var/www/satis.json /var/www/public/ -vvv > /var/log/satis.log